# docker/tool-runner-controller/supervisord.conf
[supervisord]
nodaemon=true

[program:toolrunner]
command=java -jar /usr/src/app/gaiaRunner.jar